# dist/linux/systemd.service
[Unit]
Description=Ghostty

[Service]
Type=dbus
BusName=com.mitchellh.ghostty
ExecStart=ghostty --launched-from=systemd
