# dist/linux/systemd.service.in
[Unit]
Description=Ghostty@@NAME@@

[Service]
Type=dbus
BusName=com.mitchellh.ghostty@@DEBUG@@
ExecStart=@@GHOSTTY@@ --launched-from=systemd
