# dist/linux/systemd.service.in
[Unit]
Description=@NAME@

[Service]
Type=dbus
BusName=@APPID@
ExecStart=@GHOSTTY@ --launched-from=systemd
